# sandbox/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :1 -screen 0 1024x768x24 -ac +extension RANDR
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.stdout.log
stderr_logfile=/var/log/xvfb.stderr.log

[program:xfce4]
command=startxfce4
environment=DISPLAY=":1",HOME="/root"
autorestart=true
priority=20
stdout_logfile=/var/log/xfce4.stdout.log
stderr_logfile=/var/log/xfce4.stderr.log

[program:x11vnc]
command=x11vnc -display :1 -nopw -forever -shared -rfbport 5900 -xkb
autorestart=true
priority=30
stdout_logfile=/var/log/x11vnc.stdout.log
stderr_logfile=/var/log/x11vnc.stderr.log

[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/var/log/novnc.stdout.log
stderr_logfile=/var/log/novnc.stderr.log

[program:sandbox_agent]
command=python3 -m uvicorn sandbox_agent:app --host 0.0.0.0 --port 5001
directory=/app
environment=DISPLAY=":1",WORKSPACE_DIR="/workspace"
autorestart=true
priority=50
stdout_logfile=/var/log/sandbox_agent.stdout.log
stderr_logfile=/var/log/sandbox_agent.stderr.log
